# dbraas.service
[Unit]
Description="Du bist raus." as a service

[Service]
Type=simple
ExecStart=/usr/local/bin/dbraas

[Install]
WantedBy=multi-user.target
